# misc/Linux/debian/mix-free-squidconf.prerm
#!/bin/sh
set -e

if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
	invoke-rc.d squid stop
else
	/etc/init.d/squid stop
fi
sleep 3

# restore the backup of squid configuration
cp -f /usr/share/doc/mix-free-squidconf/conf/squid.conf.orig /etc/squid/squid.conf

exit 0

