# proxytest/tags/version_00_09_017_last_with_old_protocols_support/misc/Linux/debian/mix-premium-squidconf.prerm
#!/bin/sh
set -e

if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
	invoke-rc.d squid stop || exit $?
else
	/etc/init.d/squid stop || exit $?
fi
sleep 3

# restore the backup of squid configuration
cp -f /usr/share/doc/mix-premium-squidconf/conf/squid.conf.orig /etc/squid/squid.conf
rm /etc/squid/squid-block.acl

